# crontab
*/8 * * * * php /mnt/app/public/data_collect/index.php

# always remeber the last line of a crontab file must eb an empty line
